# ssh_access_log/ssh-access-watcher.path
[Path]
PathModified=/var/log/auth.log

[Install]
WantedBy=multi-user.target
